# systemd/system/snap-core-14447.mount
[Unit]
Description=Mount unit for core, revision 14447
Before=snapd.service
After=zfs-mount.service

[Mount]
What=/var/lib/snapd/snaps/core_14447.snap
Where=/snap/core/14447
Type=squashfs
Options=nodev,ro,x-gdu.hide,x-gvfs-hide
LazyUnmount=yes

[Install]
WantedBy=default.target multi-user.target
