# systemd/system/snap-telegram\x2ddesktop-4486.mount
[Unit]
Description=Mount unit for telegram-desktop, revision 4486
Before=snapd.service
After=zfs-mount.service

[Mount]
What=/var/lib/snapd/snaps/telegram-desktop_4486.snap
Where=/snap/telegram-desktop/4486
Type=squashfs
Options=nodev,ro,x-gdu.hide,x-gvfs-hide
LazyUnmount=yes

[Install]
WantedBy=default.target multi-user.target
