# systemd/system/snap-inkscape-10509.mount
[Unit]
Description=Mount unit for inkscape, revision 10509
Before=snapd.service
After=zfs-mount.service

[Mount]
What=/var/lib/snapd/snaps/inkscape_10509.snap
Where=/snap/inkscape/10509
Type=squashfs
Options=nodev,ro,x-gdu.hide,x-gvfs-hide
LazyUnmount=yes

[Install]
WantedBy=default.target multi-user.target
